# deploy/systemd/rpi-sensor-bridge.service
[Unit]
Description=Raspberry Pi Arduino Sensor MQTT Bridge
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/rpi-sensor-bridge
EnvironmentFile=/etc/rpi-sensor-bridge/rpi-sensor-bridge.env
Environment=PYTHONPATH=/opt/rpi-sensor-bridge/src
ExecStart=/opt/rpi-sensor-bridge/.venv/bin/python -m bridge.main
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
